# supervisord.conf
[supervisord]
nodaemon=true

[program:ssserver]
command=/bin/bash -c "ssserver -c /etc/shadowsocks.json -d start"
startsecs=0